# systemd/ransomeye-git-sync.timer
# Path and File Name : /home/ransomeye/rebuild/systemd/ransomeye-git-sync.timer
# Author: nXxBku0CKFAJCBN3X1g3bQk7OxYQylg8CMw1iGsq7gU
# Details of functionality of this file: Systemd timer for automatic git sync to GitHub (runs hourly)

[Unit]
Description=RansomEye Git Sync Timer
Requires=ransomeye-git-sync.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target

