# systemd/ransomeye-retention-enforcer.timer
# Path and File Name : /home/ransomeye/rebuild/systemd/ransomeye-retention-enforcer.timer
# Author: nXxBku0CKFAJCBN3X1g3bQk7OxYQylg8CMw1iGsq7gU
# Details of functionality of this file: Systemd timer scheduling periodic runtime DB retention enforcement.

[Unit]
Description=RansomEye Runtime Retention Enforcer Timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
Unit=ransomeye-retention-enforcer.service
Persistent=true

[Install]
WantedBy=multi-user.target


